# dot_config/systemd/user/glance.service.tmpl
[Unit]
Description=Runs a Glance web server

[Service]
Type=simple
ExecStart=/home/{{ .chezmoi.username }}/.local/bin/glance -config ${HOME}/.config/glance/config.yaml

[Install]
WantedBy=default.target
